# configs/swapoff.service
[Unit]

Description=Disable swap

After=default.target

[Service]

ExecStart=/root/swapoff.sh

[Install]

WantedBy=default.target